# deployment/ansible/webhook/webhook.service
[Unit]
Description=Webhooks

[Service]
ExecStart=/usr/bin//webhook -hooks ~/webhooks/hooks.json -hotreload

[Install]
WantedBy=multi-user.target
